# deploy/systemd/garden-wake.service
[Unit]
Description=Galatea Garden 唤醒桥
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=6

[Service]
Type=simple
User=garden-wake
Group=garden-wake
WorkingDirectory=/opt/galatea-garden-wake-bridge
EnvironmentFile=/etc/galatea-garden-wake.env
ExecStart=/usr/bin/node /opt/galatea-garden-wake-bridge/dist/cli.js run
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=2
TimeoutStopSec=20
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
